# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:queue-worker]
; Queue strategy:
; - activity: audit/activity logging (high priority)
; - default: fallback for generic jobs
; Future queues:
; - notifications
; - realtime
; Scaling note:
; - php artisan queue:work redis --queue=activity
; - php artisan queue:work redis --queue=notifications
command=php artisan queue:work redis --queue=activity,default --sleep=3 --tries=3 --timeout=90
directory=/var/www/backend
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/supervisor/queue.err.log
stdout_logfile=/var/log/supervisor/queue.out.log
